# supervisor.conf
[program:nginx]
command=nginx -g "daemon off;"

[program:gunicorn]
command=/srv/env/bin/gunicorn --bin localhost:5000 hondana:app
